# sslcheck.service
[Unit]
Description=sslcheck

[Service]
WorkingDirectory=/opt/
ExecStart=/usr/bin/nodejs /opt/sslcheck/index.js

[Install]
WantedBy=multi-user.target
